# deploy/com.worldcup.review.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- Copy to ~/Library/LaunchAgents/com.worldcup.review.plist, edit paths, then bootstrap.
     Runs the daily review→re-predict→re-sim→publish loop during the tournament. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.worldcup.review</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/bot/worldcup/.venv/bin/python</string>
    <string>-m</string>
    <string>skill.helpers.cli</string>
    <string>review</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/bot/worldcup</string>
  <key>EnvironmentVariables</key>
  <dict><key>PYTHONPATH</key><string>/Users/bot/worldcup</string></dict>
  <!-- Times are LOCAL. Set them ~30 min before your match-day kickoff clusters so the
       review picks up confirmed starting XIs (lineups publish ~40 min pre-kickoff) and
       the latest scores. Example below blankets a Beijing-time tournament where WC2026
       kickoffs land 00:00-10:00 local; adjust to your timezone. Finished matches are
       skipped, so the API-Football free quota (100/day) holds even with frequent runs. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/tmp/worldcup.review.log</string>
  <key>StandardErrorPath</key><string>/tmp/worldcup.review.err</string>
</dict>
</plist>
