# deploy/com.worldcup.review.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- Copy to ~/Library/LaunchAgents/com.worldcup.review.plist, edit paths, then bootstrap.
     Runs the daily review→re-predict→re-sim→publish loop during the tournament. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.worldcup.review</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/bot/worldcup/.venv/bin/python</string>
    <string>-m</string>
    <string>skill.helpers.cli</string>
    <string>review</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/bot/worldcup</string>
  <key>EnvironmentVariables</key>
  <dict><key>PYTHONPATH</key><string>/Users/bot/worldcup</string></dict>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/tmp/worldcup.review.log</string>
  <key>StandardErrorPath</key><string>/tmp/worldcup.review.err</string>
</dict>
</plist>
